# contrib/systemd-cron.cron.weekly
#!/bin/sh
test -x /usr/libexec/systemd-cron/remove_stale_stamps || exit 0
exec /usr/libexec/systemd-cron/remove_stale_stamps
