# src/units/cron-monthly.target
[Unit]
Description=systemd-cron monthly target
StopWhenUnneeded=yes
