# src/units/cron-hourly.target
[Unit]
Description=systemd-cron hourly target
StopWhenUnneeded=yes
