# src/units/cron-schedule.timer.in
[Unit]
Description=@package@ @schedule@ timer
PartOf=cron.target
RefuseManualStart=true
RefuseManualStop=true

[Timer]
OnCalendar=@schedule@
Unit=cron-@schedule@.target

