# src/units/cron-weekly.target
[Unit]
Description=systemd-cron weekly target
StopWhenUnneeded=yes
