# src/units/cron-schedule.service.in
[Unit]
Description=@package@ @schedule@ script service
Documentation=man:systemd.cron(7)
PartOf=cron-@schedule@.target
ConditionDirectoryNotEmpty=@confdir@/cron.@schedule@
OnFailure=cron-failure@%i.service

[Service]
Type=oneshot
IgnoreSIGPIPE=false
ExecStart=@runparts@ @confdir@/cron.@schedule@
