# src/units/cron-daily.target
[Unit]
Description=systemd-cron daily target
Requires=cron-daily.service
StopWhenUnneeded=yes

