# src/units/cron-boot.timer
[Unit]
Description=systemd-cron boot timer
PartOf=cron.target
RefuseManualStart=true
RefuseManualStop=true

[Timer]
OnBootSec=0

