# src/units/cron-weekly.timer
[Unit]
Description=systemd-cron weekly timer
PartOf=cron.target
RefuseManualStart=true
RefuseManualStop=true

[Timer]
OnCalendar=weekly

