# src/units/cron-daily.target
[Unit]
Description=systemd-cron daily target
StopWhenUnneeded=yes
