# src/units/cron-boot.target
[Unit]
Description=systemd-cron boot target
StopWhenUnneeded=yes
