# src/units/cron-failure@.service.in
[Unit]
Description=@package@ OnFailure for %i
Documentation=man:systemd.cron(7)

[Service]
Type=oneshot
ExecStart=@libdir@/systemd-cron/mail-on-failure.sh %i
