# src/units/cron-mail@.service.in
[Unit]
Description=systemd-cron job mail — %i
Documentation=man:systemd.cron(7)
# TODO: drop conditional and roll back to straight OnSuccess= if we ever bump past the systemd ≥ 236 requirement (this is 249); cf. #165
@remove_if_no_onsuccess@RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
PassEnvironment=SENDMAIL
ExecStart=@libexecdir@/systemd-cron/mail_for_job %i
User=_cron-failure
# TODO: drop Group= systemd ≥ 236 requirement (this is 245); cf. https://github.com/systemd-cron/systemd-cron/issues/165#issuecomment-2661583788
Group=systemd-journal
