# src/units/cron-schedule.target.in
[Unit]
Description=@package@ @schedule@ target
Requires=cron-@schedule@.service
StopWhenUnneeded=yes

