# src/units/cron-monthly.timer
[Unit]
Description=Sdcron monthly timer
PartOf=cron.target
RefuseManualStart=true
RefuseManualStop=true

[Timer]
OnCalendar=monthly

