# src/units/cron.target.in
[Unit]
Description=@package@
@requires@

[Install]
WantedBy=multi-user.target

