# src/units/cron-weekly.target
[Unit]
Description=systemd-cron weekly target
Requires=cron-weekly.service
StopWhenUnneeded=yes

