# src/units/cron-update.path.in
[Unit]
Description=@package@ path monitor
Documentation=man:systemd.cron(7)

[Path]
PathChanged=/etc/crontab
PathChanged=/etc/cron.d
PathChanged=/etc/anacrontab
PathChanged=@statedir@
