# src/units/cron.target.in
[Unit]
Description=@package@
Documentation=man:systemd.cron(7)
@requires@

[Install]
WantedBy=multi-user.target

