# src/units/cron-boot.service.in
[Unit]
Description=systemd-cron @schedule@ script service
Documentation=man:systemd.cron(7)
PartOf=cron-@schedule@.target
ConditionDirectoryNotEmpty=/etc/cron.@schedule@
ConditionPathExists=!/run/crond.bootdir
OnFailure=cron-failure@%i.service

[Service]
Type=oneshot
IgnoreSIGPIPE=false
ExecStart=@runparts@ /etc/cron.@schedule@
ExecStartPost=/usr/bin/touch /run/crond.bootdir
