# src/units/cron.target.in
[Unit]
Description=systemd-cron
Documentation=man:systemd.cron(7)
@requires@
Wants=cron-update.path

[Install]
WantedBy=multi-user.target

