# src/units/cron-monthly.target
[Unit]
Description=systemd-cron monthly target
Requires=cron-monthly.service
StopWhenUnneeded=yes

