# src/units/cron-schedule.target.in
[Unit]
Description=@package@ @schedule@ target
Documentation=man:systemd.cron(7)
Requires=cron-@schedule@.service
StopWhenUnneeded=yes

