# src/units/cron-boot.target
[Unit]
Description=systemd-cron boot target
Requires=cron-boot.service
StopWhenUnneeded=yes

