# src/units/cron-hourly.target
[Unit]
Description=systemd-cron hourly target
Requires=cron-hourly.service
StopWhenUnneeded=yes

