# src/units/cron-update.path.in
[Unit]
Description=@package@ path monitor
Documentation=man:systemd.cron(7)
RefuseManualStart=true
RefuseManualStop=true
After=basic.target
Before=cron.target

[Path]
PathChanged=/etc/crontab
PathChanged=/etc/cron.d
PathChanged=/etc/anacrontab
PathChanged=/var/spool/cron
