# src/units/cron-schedule.service.in
[Unit]
Description=systemd-cron @schedule@ script service
Documentation=man:systemd.cron(7)
PartOf=cron-@schedule@.target
ConditionDirectoryNotEmpty=/etc/cron.@schedule@
OnFailure=cron-failure@%i.service

[Service]
Type=oneshot
IgnoreSIGPIPE=false
LogLevelMax=@use_loglevelmax@
ExecStart=@runparts@ /etc/cron.@schedule@
