# src/units/cron-failure@.service.in
[Unit]
Description=@package@ OnFailure for %i
Documentation=man:systemd.cron(7)
RefuseManualStart=true
RefuseManualStop=true
ConditionFileIsExecutable=/usr/sbin/sendmail

[Service]
Type=oneshot
ExecStart=@libdir@/@package@/mail_on_failure %i
DynamicUser=yes
Group=systemd-journal
