# src/units/cron-update.service.in
[Unit]
Description=@package@ update units
Documentation=man:systemd.cron(7)

[Service]
Type=oneshot
ExecStartPre=/usr/bin/touch /run/crond.reboot
ExecStart=/bin/systemctl daemon-reload
ExecStartPost=/bin/systemctl try-restart cron.target
