# src/units/cron-after-var.service.in
[Unit]
Description=@package@ Rerun generator only if /var is a separate mount
Documentation=man:systemd.cron(7)
After=local-fs.target
ConditionPathExists=!@statedir@/.ok

[Service]
Type=oneshot
ExecStartPre=/bin/rm -f /run/crond.reboot
ExecStart=/bin/systemctl daemon-reload
ExecStartPost=/bin/systemctl try-restart cron.target
