# src/units/cron-boot.timer.in
[Unit]
Description=@package@ @schedule@ timer
Documentation=man:systemd.cron(7)
PartOf=cron.target

[Timer]
OnBootSec=0
Unit=cron-@schedule@.target
