# src/units/cron-boot.timer.in
[Unit]
Description=@package@ @schedule@ timer
PartOf=cron.target
RefuseManualStart=true
RefuseManualStop=true

[Timer]
OnBootSec=0
Unit=cron-@schedule@.target

