# src/units/cron-update.service.in
[Unit]
Description=systemd-cron update units
Documentation=man:systemd.cron(7)

[Service]
Type=oneshot
ExecStartPre=/usr/bin/touch /run/crond.reboot
ExecStart=/bin/sh -c 'systemctl daemon-reload ; sleep 1 ; systemctl restart cron.target'
